# roles/h2o/files/etc/init/h2o.conf
start on runlevel [2345]
stop on runlevel [016]

exec /usr/local/h2o/bin/h2o -c /etc/h2o/h2o.yml

respawn
